# supervisord.conf
; Security Bundle - supervisord configuration
; Manages all processes in the all-in-one container

[supervisord]
nodaemon=true
logfile=/var/log/security-bundle/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/security-bundle

; ── Nginx + ModSecurity ─────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/security-bundle/nginx-stdout.log
stderr_logfile=/var/log/security-bundle/nginx-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── HA Controller (Consul leader election + VIP failover) ────────────
[program:hactl]
command=/usr/local/bin/hactl --config /config/ha-controller.yml --node %(ENV_HACTL_NODE_ID)s --log-level info
autostart=%(ENV_HACTL_ENABLED)s
autorestart=true
priority=15
startsecs=5
stdout_logfile=/var/log/security-bundle/hactl-stdout.log
stderr_logfile=/var/log/security-bundle/hactl-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── CrowdSec ───────────────────────────────────────────────────────────
[program:crowdsec]
command=/usr/bin/crowdsec -c /etc/crowdsec/config.yaml
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/security-bundle/crowdsec-stdout.log
stderr_logfile=/var/log/security-bundle/crowdsec-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Fail2Ban ───────────────────────────────────────────────────────────
[program:fail2ban]
command=/usr/bin/fail2ban-server -f -x -v start
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/security-bundle/fail2ban-stdout.log
stderr_logfile=/var/log/security-bundle/fail2ban-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Wazuh Manager ──────────────────────────────────────────────────────
; wazuh-control has no foreground mode; wrapper starts it and monitors analysisd
[program:wazuh-manager]
command=/opt/security-bundle/scripts/wazuh-manager-wrapper.sh
autostart=true
autorestart=true
priority=40
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/var/log/security-bundle/wazuh-manager-stdout.log
stderr_logfile=/var/log/security-bundle/wazuh-manager-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Wazuh Indexer (OpenSearch) ─────────────────────────────────────────
; Do NOT pass -Epath.home as the package already sets it via jvm.options/env
[program:wazuh-indexer]
command=/usr/share/wazuh-indexer/bin/opensearch
user=wazuh-indexer
autostart=true
autorestart=true
priority=35
startsecs=30
environment=OPENSEARCH_JAVA_OPTS="-Xms384m -Xmx384m",OPENSEARCH_PATH_CONF="/etc/wazuh-indexer"
stdout_logfile=/var/log/security-bundle/wazuh-indexer-stdout.log
stderr_logfile=/var/log/security-bundle/wazuh-indexer-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Wazuh Dashboard ────────────────────────────────────────────────────
[program:wazuh-dashboard]
command=/usr/share/wazuh-dashboard/bin/opensearch-dashboards --config /etc/wazuh-dashboard/opensearch_dashboards.yml
user=wazuh-dashboard
autostart=true
autorestart=true
priority=50
startsecs=30
environment=NODE_OPTIONS="--max-old-space-size=256"
stdout_logfile=/var/log/security-bundle/wazuh-dashboard-stdout.log
stderr_logfile=/var/log/security-bundle/wazuh-dashboard-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Watchdog (health monitoring + nginx failover) ──────────────────────
[program:watchdog]
command=/opt/security-bundle/scripts/watchdog.sh
autostart=true
autorestart=true
priority=60
stdout_logfile=/var/log/security-bundle/watchdog-stdout.log
stderr_logfile=/var/log/security-bundle/watchdog-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Log Shipper (tails alerts -> MySQL) ────────────────────────────────
[program:log-shipper]
command=/usr/bin/python3 /opt/security-bundle/scripts/log-shipper.py
autostart=true
autorestart=true
priority=70
stdout_logfile=/var/log/security-bundle/log-shipper-stdout.log
stderr_logfile=/var/log/security-bundle/log-shipper-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
